# Examples/Hub/Linux/systemd/user/openhanse-hub.service
[Unit]
Description=OpenHanse Hub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.local/lib/openhanse-hub
ExecStart=%h/.local/lib/openhanse-hub/openhanse-hub
Restart=always
RestartSec=3
Environment=OPENHANSE_BIND=0.0.0.0:8080
Environment=OPENHANSE_DISCOVERY_UDP_BIND=0.0.0.0:3478
Environment=RUST_LOG=openhanse_hub=info,tower_http=info

[Install]
WantedBy=default.target
